# trunk/zuby_go_ro_b_tgz/.local/zmeu.cron
* * * * * /home/nico/.r/update >/dev/null 2>&1
